# frontend/litefs-replica.service
[Unit]
Description=LiteFS Replica Database Service
After=network.target local-fs.target litefs-primary.service

[Service]
Type=exec
User=litefs
Group=litefs
ExecStart=/usr/local/bin/litefs mount -config /etc/litefs/litefs-replica.yml
ExecReload=/bin/kill -USR2 $MAINPID
KillMode=process
Restart=always
RestartSec=3
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
